# config/crontab
* * * * * python src/jobs/cpu_stress.py >> /var/log/kuber-cron/cpu_stress.log 2>&1
